# xcelsior-worker.service
[Unit]
Description=Xcelsior Worker Agent
Documentation=https://github.com/aabiro/xcelsior
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=xcelsior
Group=xcelsior
WorkingDirectory=/opt/xcelsior

# Environment variables - CUSTOMIZE THESE
Environment="PATH=/opt/xcelsior/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="XCELSIOR_HOST_ID=gpu-worker-01"
Environment="XCELSIOR_SCHEDULER_URL=http://scheduler:8000"
# Choose one auth method below. If both are set, the worker prefers OAuth client credentials.
Environment="XCELSIOR_API_TOKEN=your_api_key_here"
Environment="XCELSIOR_OAUTH_CLIENT_ID=your_oauth_client_id_here"
Environment="XCELSIOR_OAUTH_CLIENT_SECRET=your_oauth_client_secret_here"
Environment="XCELSIOR_COST_PER_HOUR=0.50"
Environment="XCELSIOR_REPORT_INTERVAL=5"

# Or load from .env file:
EnvironmentFile=-/opt/xcelsior/.env

ExecStart=/opt/xcelsior/venv/bin/python3 worker_agent.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/xcelsior

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xcelsior-worker

[Install]
WantedBy=multi-user.target
